# com.pharmadigest.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pharmadigest.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jiagabu/pharma-digest/run.sh</string>
  </array>

  <!-- Run daily at 07:10 Taipei time (macOS local time must be set to Asia/Taipei) -->
  <!-- Also check every hour as fallback if Mac was asleep at 07:10 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jiagabu/pharma-digest</string>

  <key>StandardOutPath</key>
  <string>/Users/jiagabu/pharma-digest/logs/launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jiagabu/pharma-digest/logs/launchd_error.log</string>

  <!-- Restart if it crashes -->
  <key>KeepAlive</key>
  <false/>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
